# supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod --fork --logpath /var/log/mongodb/mongod.log

[program:java-app]
command=/usr/bin/java -jar /url-shortener-images.jar